# supervisor/conf.d/uwsgi.conf
[program:uwsgi]
command=/usr/local/bin/uwsgi
	--master --socket localhost:8001 --workers=1 --module timetracker.wsgi --home /home/ubuntu/ansrportal/ANSR_ENV
autorestart=true
startretries=1